# chef-repo/cookbooks/storm/files/default/storm-ui.conf
[program:storm-ui]
command=/opt/storm/install/storm/current/bin/storm ui
autostart=true
autorestart=true
user=storm
redirect_stderr=true
stdout_logfile=/opt/storm/install/storm/log/storm-ui.log
